# 81.deploy-django-js-project-with-nginx-web-server-via-github-and-automate-deployment/django-project.com.gunicorn.socket
[Unit]
Description=django-project.com.gunicorn socket

[Socket]
ListenStream=/run/django-project.com.gunicorn.sock

[Install]
WantedBy=sockets.target
